# com.ayan.engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE PLIST PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ayan.engine</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/ayan-engine</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/shoully/Library/Logs/ayan-engine.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shoully/Library/Logs/ayan-engine.err.log</string>
</dict>
</plist>
